# debian/raumzeitchef-ircbot.service
[Unit]
Description=RaumZeitChef IRC Bot

[Service]
User=raumzeitchef
ExecStart=/usr/bin/ircbot-chef 

[Install]
WantedBy=multi-user.target
